# debian/shellfleet-agent.service
[Unit]
Description=ShellFleet Agent
After=network-online.target docker.service k3s.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root

ExecStart=/usr/bin/shellfleet-agent

# State (auth token) lives at /etc/shellfleet/agent-token.txt; read-only fields
# go in /etc/shellfleet/env. A bind mount of /etc lets the config editor see
# host configs even on systems with strict ProtectHome/ProtectSystem defaults.
EnvironmentFile=-/etc/shellfleet/env

Restart=on-failure
RestartSec=5s

# Modest hardening — the agent still needs root for the current systemd,
# package, Docker, and config surfaces, but may not gain additional privileges
# through setuid/setgid binaries or file capabilities.
NoNewPrivileges=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
